# deploy/noc-webhook.service
[Unit]
Description=NOC Deploy Webhook Receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/app/deploy
ExecStart=/usr/bin/node /home/deploy/app/deploy/webhook-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
